# config/supervisord.ini
[program:license_web]
command=rails s --binding=0.0.0.0           ; the program (relative uses PATH, can take args)
directory=/srv/license_website              ; directory to cwd to before exec (def no cwd)